# deploy/files/carto.service
[Unit]
Description=CartoDB

[Service]
Type=simple
ExecStart=/root/cartodb/init.sh

[Install]
WantedBy=multi-user.target